# .docker/unit/services.d/unitd/run
#!/bin/sh

exec unitd --no-daemon --user root --group root --log /var/log/nginx-unit/nginx-unit.log
