# .docker/unit/services.d/pcitd-agent/finish
#!/bin/sh

set -x

s6-svscanctl -t /etc/services.d
