# .docker/unit/services.d/pcitd-agent/run
#!/bin/sh

set -x

cd /app/pcit

exec /app/pcit/bin/pcitd agent
