# .docker/unit/services.d/redis-server/run
#!/bin/sh

set -x

cd /data

exec /usr/local/bin/redis-server /usr/local/etc/redis/redis.conf
