# .docker/unit/services.d/dockerd/run
#!/bin/sh

rm -rf /var/run/docker.pid

dockerSocket='unix:///var/run/docker.sock'

/usr/local/bin/dind dockerd --host="$dockerSocket" --host=tcp://127.0.0.1:2375
