# .docker/unit/services.d/pcitd/run
#!/bin/sh

set -x

cd /app/pcit

exec /app/pcit/bin/pcitd ${PCITD_CMD:-server}
